# kauth.service
[Unit]
Description=kauth identity provider
After=network.target

[Service]
Type=simple
User=lars
EnvironmentFile=/home/lars/kauth-go/.env
WorkingDirectory=/home/lars/kauth-go
ExecStart=/home/lars/kauth-go/bin/kauth-arm64
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/lars/kauth-go/data /home/lars/kauth-go/static

# Minnebegrensning: GOMEMLIMIT lar Go GC-e seg ned før cgroup-taket nås,
# slik at MemoryMax sjelden utløser OOM-kill. Tomgang ~21 MB, last ~30-50 MB.
Environment=GOMEMLIMIT=96MiB
MemoryHigh=112M
MemoryMax=128M

[Install]
WantedBy=multi-user.target
